# root/etc/sv/swayRunDir/run
#!/usr/bin/env sh

exec 2>&1

ID="1000"
RUNTIME_DIR="/run/user/$ID"

if [ ! -d $RUNTIME_DIR ];
then
	mkdir -p $RUNTIME_DIR
	chown $ID $RUNTIME_DIR -R;
	chmod 0700 $RUNTIME_DIR -R;
fi
